# example/etc/crontabs/pull
# min   hour    day     month   weekday   user    command
*/1     *       *       *       *         pull    /home/pull/rsync-pull
